# system_files/etc/profile.d/filotimo-ime.sh
export INPUT_METHOD=fcitx
export QT_IM_MODULES="wayland;fcitx;ibus"
export GLFW_IM_MODULE=ibus
export XMODIFIERS=@im=fcitx
export GTK_IM_MODULE=fcitx
